# summariz.service
[Unit]
Description=YouTube Content Summarizer API
After=network.target

[Service]
Type=simple
User=sajja
Group=sajja
WorkingDirectory=/home/sajja/summariz
Environment="PATH=/home/sajja/summariz/venv/bin:/usr/local/bin:/usr/bin:/bin"
Environment="HOME=/home/sajja"
EnvironmentFile=/home/sajja/summariz/.env
ExecStart=/home/sajja/summariz/venv/bin/uvicorn app:app --host 127.0.0.1 --port 8000
Restart=always
RestartSec=5

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=summariz

ProtectHome=false

[Install]
WantedBy=multi-user.target
